# grafana_agent/templates/etc_rc.d_grafana-agent
#!/bin/sh

# PROVIDE: grafana-agent
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable grafana agent
# grafana_agent_enable="YES"

. /etc/rc.subr

name=grafana_agent
rcvar=grafana_agent_enable

load_rc_config $name

: ${grafana_agent_enable:="NO"}
: ${grafana_agent_config:="{{ grafana_agent_config_file }}"}

pidfile="/var/run/${name}.pid"
required_files="${grafana_agent_config}"

procname="{{ grafana_agent_binary_file }}"
command="/usr/sbin/daemon"
command_args="-p ${pidfile} -t ${name} -o {{ grafana_agent_log_file }} \
  /usr/bin/env ${grafana_env} ${procname} \
  -server.http.address=127.0.0.1:8080 \
  -server.grpc.address=127.0.0.1:9095 \
  -config.file=${grafana_agent_config}"

run_rc_command "$1"
